# scripts/inittab
::sysinit:/etc/init.d/rcS
::respawn:/bin/busybox cttyhack /etc/shell
::ctrlaltdel:/sbin/reboot -f
::shutdown:/sbin/poweroff -f
