# scripts/inittab
null::sysinit:/sbin/ip link set lo up
null::sysinit:/etc/init.d/rcS
tty1::respawn:/bin/busybox cttyhack /etc/shell
::ctrlaltdel:/sbin/reboot -f
::shutdown:/sbin/poweroff -f
